# templates/com.scaviefae.simpleloop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scaviefae.simpleloop.{{PROJECT_NAME}}</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{LIB_DIR}}/daemon.sh</string>
        <string>{{PROJECT_DIR}}</string>
        <string>{{INTERVAL}}</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{PROJECT_DIR}}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>{{LOG_FILE}}</string>

    <key>StandardErrorPath</key>
    <string>{{LOG_FILE}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>{{PATH}}</string>
        <key>GIT_MAIN_BRANCH</key>
        <string>{{GIT_MAIN_BRANCH}}</string>
{{LOOP_LANE_KEY}}
{{LOOP_LANE_VAL}}
    </dict>
</dict>
</plist>
